# .config/profile.d/000-path.sh
export PATH=$PATH:/sbin:/usr/sbin/:${HOME}/bin:${HOME}/opt/install/bin
