# .config/profile.d/terragrunt.sh
alias tg=terragrunt
alias tga='terragrunt apply'
alias tgp='terragrunt plan'
alias tgpu='terragrunt plan --terragrunt-source-update'
alias tgs='terragrunt show'
